# files/server/bin/teamcity-server
#!/bin/sh

### BEGIN INIT INFO
# Provides:          teamcity-server
# Required-Start:    $mysql
# Required-Stop:     $mysql
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start and stop the TeamCity Build Server
### END INIT INFO

# For CentOS/RHEL
#chkconfig: 2345 80 05
#description: Start and stop the TeamCity Build Server

TEAMCITY_HOME=/opt/teamcity-server
TEAMCITY_USER=teamcity
SERVER_SCRIPT=${TEAMCITY_HOME}/bin/server.sh
SERVER_LOG=/var/log/teamcity/server-console.log

umask 002

case "$1" in
  start)
    su - ${TEAMCITY_USER} -c "${SERVER_SCRIPT} $1 > ${SERVER_LOG} 2>&1"
    RETVAL=$?
    ;;

  stop)
    su - ${TEAMCITY_USER} -c "${SERVER_SCRIPT} $1 >> ${SERVER_LOG} 2>&1"
    RETVAL=$?
    ;;

  restart)
    $0 stop && $0 start
    RETVAL=$?
    ;;

  *)
    echo "Usage: $0 { start | stop | restart }"
    exit 1
    ;;

esac
exit ${RETVAL}

